# supervisor.d/service.ini
[supervisord]
nodaemon=true
user=root

[program:crond]
command         = /usr/sbin/crond -f
directory       = /app
autostart       = true
autorestart     = true
startsecs       = 5
numprocs        = 1
stopsignal      = INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr = true

[program:mosdns]
command         = /usr/bin/mosdns start --dir /app/mosdns
directory       = /app/mosdns
autostart       = true
autorestart     = true
startsecs       = 5
numprocs        = 1
stopsignal      = INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr = true
